# dist/lib/systemd/system/fxpixel.service
[Unit]
Description=fxPixel service

[Service]
User=root
Group=root
ExecStart=/usr/bin/fxpixel start --config=/etc/fxpixel/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
