# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable fxpixel.service
systemctl start fxpixel.service
