# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop fxpixel.service
systemctl disable fxpixel.service
